# script.cron
0 */1 * * * cd ~/AASF-hacktoberfest-api/src && /usr/bin/node updaterDB.js >> logs.txt 2>&1