# script.cron
19 */1 * * * cd ~/git-init-FOSS-API/src && /usr/bin/node updaterDB.js >> logs.txt 2>&1